# deploy/trading-bot.service
[Unit]
Description=Trading bot (daemon mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/trading-bot
EnvironmentFile=/opt/trading-bot/.env
ExecStart=/usr/bin/npx tsx index.ts BTCUSDT 5m --daemon
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
